# etc/systemd/system/apply-monitor-rules.service
[Unit]
Description=Applica regole firewall monitoraggio al boot
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=oneshot
ExecStartPre=/bin/sh -c 'until ping -c1 -W1 192.168.20.1 >/dev/null 2>&1; do sleep 2; done'
ExecStart=/usr/local/bin/apply-monitor-rules.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
