# etc/systemd/system/check-traffic.timer
[Unit]
Description=Esegue il controllo traffico ogni 5 minuti

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
